# supervisor.d/rq_minty_all_queues.conf
[program:rq_minty_all_queues]
; Point the command to the specific rq command you want to run.
; If you use virtualenv, be sure to point it to
; /path/to/virtualenv/bin/rq
; Also, you probably want to include a settings module to configure this
; worker.  For more info on that, see http://python-rq.org/docs/workers/

command=/home/mintviz/minty/minty_env/bin/rq worker --url redis://127.0.0.1:16379/1 --name %(program_name)s-%(process_num)s --path /minty high low

user=mintviz

environment=MINTCAST_PATH="/home/mintviz/mintcast",MINTCAST_IS_ON_SERVER="NO",PATH="/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/home/mintviz/.local/bin:/home/mintviz/bin"

redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s-%(process_num)s.log

; process_num is required if you specify >1 numprocs
process_name=%(program_name)s-%(process_num)s
numprocs=2

enable-threads=true

directory=/home/mintviz/minty

stopsignal=TERM

; These are up to you
autostart=true
autorestart=true
